# mnt-nfs.automount
[Unit]
Description=Automount Nfs

[Automount]
Where=/mnt/nfs

[Install]
WantedBy=multi-user.target
